# templates/mockserver.service.erb
[Unit]
Description=Mockserver
After=network.target

[Service]
User=<%= @user %>
ExecStart=<%= @dir %>/bin/mockserver.sh start
ExecStop=<%= @dir %>/bin/mockserver.sh stop

[Install]
WantedBy=multi-user.target
