# templates/mockserver.service.erb
[Unit]
Description=Mockserver
After=network.target

[Service]
User=<%= scope.lookupvar('mockserver::user') %>
ExecStart=<%= scope.lookupvar('mockserver::dir') %>/bin/mockserver.sh start
ExecStop=<%= scope.lookupvar('mockserver::dir') %>/bin/mockserver.sh stop

[Install]
WantedBy=multi-user.target
